# deployments/claude-ops.service
[Unit]
Description=Claude Ops — GitHub issue to PR automation
After=network.target
Documentation=https://github.com/gs97ahn/claude-ops

[Service]
Type=simple
User=<REPLACE_WITH_OPERATOR_USER>
Group=<REPLACE_WITH_OPERATOR_GROUP>
WorkingDirectory=/srv/claude-ops
EnvironmentFile=/etc/claude-ops/.env
ExecStart=/usr/local/bin/claude-ops -config /etc/claude-ops/config.yaml
Restart=on-failure
RestartSec=5s

# Ensure child claude processes are killed on stop
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/srv/claude-ops/data
ReadWritePaths=/srv/claude-ops/.worktrees

[Install]
WantedBy=multi-user.target
